# conf/units/fatbuildrd.service
[Unit]
Description=Fatbuildr build daemon

[Service]
BusName=org.rackslab.Fatbuildr
ExecStart=/usr/bin/fatbuildrd
WatchdogSec=3min
Restart=on-failure

[Install]
WantedBy=multi-user.target
