# conf/units/fatbuildrd.path
[Path]
DirectoryNotEmpty=/var/lib/fatbuildr/queue

[Install]
WantedBy=multi-user.target
